# backend/systemd/overcyber.service.example
[Unit]
Description=Overcyber backend (axum + sqlite)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/overcyber-backend
Restart=on-failure
RestartSec=3

DynamicUser=yes
StateDirectory=overcyber
RuntimeDirectory=overcyber
Environment="DATABASE_PATH=/var/lib/overcyber/overcyber.db"
Environment="BIND_ADDR=127.0.0.1:8787"
Environment="PUBLIC_ORIGIN=https://overcyber.example.com"
Environment="RUST_LOG=info,overcyber_backend=info"
# EnvironmentFile=/etc/overcyber/overcyber.env   # contém SESSION_SECRET e ADMIN_BOOTSTRAP_PASSWORD

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
